# trunk/floater/debian/floater-server.cron.d
0 3 * * Sun	su floater /usr/lib/floater/res2html > /dev/null
